# conf/lickur.conf
[program:lickur]
command = /usr/bin/python /srv/Lickur/run.py server
priority = 1
